# rpm/fpm/preremove
#!/bin/bash

set +e

systemctl stop php-fpm-$VERSION.service
systemctl disable php-fpm-$VERSION.service

exit 0