# rpm/preremove
#!/bin/bash

set +e

systemctl stop php-fpm.service
systemctl disable php-fpm.service

exit 0