# alpine/fpm/pre-deinstall
#!/bin/bash

set +e

rc-update php-fpm-$VERSION stop
rc-update del php-fpm-$VERSION default

exit 0