# alpine/fpm/pre-install
#!/bin/bash

set +e
VERSION=

WWWDATAUSER=$(id www-data)
if [ $? -eq 1 ]; then
    echo "Creating www-data user"
    adduser -S www-data
fi

rc-service php-fpm-$VERSION stop
rc-update del php-fpm-$VERSION default

$(which mkdir) -p /var/run/php/$VERSION
$(which touch) /var/run/php/$VERSION/php-fpm.pid

exit 0