# alpine/fpm/pre-deinstall
#!/bin/bash

set +e

VERSION=

rc-update php-fpm-$VERSION stop
rc-update del php-fpm-$VERSION default

exit 0